# common/atlassian-cloud-backuper.cron
################################################################################

SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/srv/atlassian-backups

################################################################################

# Backup Jira data
0 3 * * Sun root atlassian-cloud-backuper jira &> /dev/null || :

# Backup Confluence data
0 9 * * Sun root atlassian-cloud-backuper confluence &> /dev/null || :
